# java_console/bin/etc/systemd/system/gerefi.service
#
# sudo systemctl enable gerefi
#

[Unit]
Description=gerefi systemd

[Service]
ExecStart=/rw/startup_root.sh

[Install]
WantedBy=multi-user.target
